# srcs/mysql/supervisord.conf
[supervisord]
nodaemon = true
user=root

[program:mysql]
command=sh -c "mysqld && kill -s SIGINT `cat supervisord.pid`"
