# srcs/ftps/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:ftps]
command=sh -c "vsftpd /etc/vsftpd/vsftpd.conf && kill -s SIGINT `cat supervisord.pid`"
